# systemd/appdata-backup.timer
[Unit]
Description=Schedule appdata backups

[Timer]
OnCalendar=04:10
Persistent=true

[Install]
WantedBy=multi-user.target
