# systemd/b2-sync.service
[Unit]
Description=Sync Backups to B2

[Service]
Type=oneshot

ExecStart=/usr/bin/env sh -c '/home/sitheris/docker/bin/b2_sync.sh'
